# deploy/pkg/virtual-flex-rigctld.service
[Unit]
Description=Hamlib rigctld for virtual-flex
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/virtual-flex/rigctld.env
ExecStart=/usr/bin/rigctld -m ${RIG_MODEL} -r ${K4_HOST}:${K4_CAT_PORT} -t ${RIGCTLD_PORT}
Restart=always
RestartSec=5
DynamicUser=yes

[Install]
WantedBy=multi-user.target
