# deploy/pkg/virtual-flex.service
[Unit]
Description=virtual-flex (native Elecraft K4 -> 4O3A Genius stack)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/usr/lib/virtual-flex
ExecStart=/usr/bin/python3 -m virtualflex --config /etc/virtual-flex/config.toml --log-level INFO
Restart=always
RestartSec=5
DynamicUser=yes

[Install]
WantedBy=multi-user.target
